# hardware/common/etc/dln-netmon@.service
[Unit]
Description=Monitor dln network hardware

[Service]
Type=oneshot
ExecStart=/opt/dlt/bin/net_meas /var/log/dlt/measurements %i
