# hardware/common/etc/dln-sysmon.service
[Unit]
Description=Monitor dln system hardware

[Service]
Type=oneshot
ExecStart=/opt/dlt/bin/sys_meas /var/log/dlt/measurements
