# hardware/base/etc/network/if-pre-up.d/iptables
#!/bin/sh

PATH=/sbin

# flush existing rules
iptables -F
iptables -F -t nat
iptables -X

# NAT for wdln ferries
iptables -t nat -I POSTROUTING -o eth0 -j MASQUERADE
